# com.emailagent.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.emailagent.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ramasamymolaiappan/Projects/repo/ramasamy-personal-agent/.venv/bin/python3</string>
        <string>/Users/ramasamymolaiappan/Projects/repo/ramasamy-personal-agent/agent.py</string>
    </array>

    <!-- Run every day at 8:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ramasamymolaiappan/Projects/repo/ramasamy-personal-agent/logs/agent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ramasamymolaiappan/Projects/repo/ramasamy-personal-agent/logs/agent.error.log</string>

    <!-- Run when Mac wakes if the scheduled time was missed -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
